# debian/wazo-backup.cron.d
PATH=/usr/sbin

5 0 * * *   root       wazo-backup data /var/backups/wazo/data ; wazo-backup db /var/backups/wazo/db
