# debian/xivo-backup.cron.d
PATH=/usr/sbin

5 0 * * *   root       wazo-backup data /var/backups/xivo/data ; wazo-backup db /var/backups/xivo/db
