# debian/xivo-backup.cron.d
PATH=/usr/sbin

5 0 * * *   root       xivo-backup data /var/backups/xivo/data ; xivo-backup db /var/backups/xivo/db
